# docker/instack-virt/libvirtd.ini
[program:libvirtd]
command=/usr/sbin/libvirtd
